# docs/server/etc/systemd/system/humo_bot.socket
[Unit]
Description=humo_bot socket

[Socket]
ListenStream=/run/humo_bot.sock

[Install]
WantedBy=sockets.target
